# deploy/systemd/nextapp.service
[Unit]
Description=Thermovault Next.js App
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/your/app   # REPLACE
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=thermovault-next

[Install]
WantedBy=multi-user.target
